# .config/systemd/user/mail-sync@.target
[Unit]
PartOf=%p.target

[Install]
WantedBy=%p.target
