# .config/systemd/user/mbsync.service
[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -c %E/mbsync/config -V %i
Restart=on-failure
RestartSec=1min
