# .config/systemd/user/msmtp.path
[Path]
DirectoryNotEmpty=%h/.local/share/mail_queue
MakeDirectory=true
DirectoryMode=0700
